# crontab
0 * * * * pacman -Syu 
0 */3 * * * pacman -Sc
