# ansible/files/crontab
@reboot sleep 15 && DISPLAY=:0 /usr/local/bin/profile-selector
